# system/etc/NetworkManager/dispatcher.d/00-indicator-update.sh
#!/bin/bash

pkill -RTMIN+3 i3blocks

exit 0

